# .config/systemd/user/change-wallpaper.timer.d/override.conf
[Unit]
Description=Change the wallpaper on X display :0 every hour
Requires=change-wallpaper.service

[Timer]
OnBootSec=1second
OnUnitActiveSec=1minute

[Install]
WantedBy=timers.target
